# etc/systemd/system/home-liolin-shares-ds.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Before=remote-fs.target
After=network-online.target
Wants=network-online.target

[Mount]
Where=/home/liolin/shares/ds
What=//192.168.1.176/data
Type=cifs
Options=vers=2.0,credentials=/home/liolin/.credentials/data,uid=1000,gid=1000,iocharset=utf8

[Install]
WantedBy=multi-user.target
